# dot_config/systemd/user/downloads-cleanup.service
[Unit]
Description=Clean old files from Downloads directory

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-tmpfiles --user --clean %h/.config/user-tmpfiles.d/downloads-cleanup.conf
